# etc/xdg/autostart/deckscale.desktop
[Desktop Entry]
Comment=Scale the desktop for handheld devices
Exec=/usr/bin/deckscale
Name=Steam Deck desktop rescaler
Type=Application
