# etc/xdg/autostart/steam.desktop
[Desktop Entry]
Comment=Start Steam sliently
Exec=/usr/bin/steam -silent -steamdeck
Name=Steam
Type=Application
